# filesystem/etc/system/pidora-keyboard.service
[Unit]
After=default.target

[Service]
ExecStart=/root/pidora-keyboard/pidora-keyboard.py start

[Install]
WantedBy=default.target